# liveenv/root/etc/profile.d/alias.sh
alias ls='ls -lhF --color'
alias vi='vim'
